# ffpb/ffpb-openwrt-migration-tools/files/etc/uci-defaults/90-tlwr1043ndv2-autoupdater-to-stable
#!/bin/sh

[ "$(cat /tmp/sysinfo/board_name)" == "tl-wr1043nd-v2" ] && \
  uci set autoupdater.settings.branch="stable" && \
  uci set autoupdater.settings.enabled=1

exit 0
