# ffho/ffho-openwrt-migration-tools/files/etc/uci-defaults/90-tl-841v10-autoupdater-to-stable
#!/bin/sh
[ "$(cat /tmp/sysinfo/model)" == "TP-Link TL-WR841N/ND v10" ] && \
	uci set autoupdater.settings.branch="stable" && \
	uci set autoupdater.settings.enabled=1
exit 0
